# runpod/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:stt]
command=/opt/speak2/run-stt.sh
directory=/opt/moshi
autostart=true
autorestart=true
startretries=10
startsecs=120
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=100

[program:tts]
command=/opt/speak2/run-tts.sh
directory=/opt/moshi
autostart=true
autorestart=true
startretries=10
startsecs=120
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=101

[program:backend]
command=/opt/speak2/run-backend.sh
directory=/opt/unmute
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=200

[program:frontend]
command=/opt/speak2/run-frontend.sh
directory=/opt/frontend
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=300

[program:traefik]
command=/usr/local/bin/traefik --configFile=/etc/traefik/traefik-static.yaml
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=400
